# decommission-tool-deb/src/deb/init/decommission-tool.conf
description "Decommission Tool REST Service"
author "Tim Squires"

start on runlevel [2345]
stop on runlevel [!2345]

setuid decommission
chdir /opt/decommission-tool

script
    ./run.sh
end script
