# decommission-tool-deb/src/deb/init/healthcheck.conf
description "Health check REST Service"
author "Tim Squires"

start on runlevel [2345]
stop on runlevel [!2345]

setuid healthcheck
chdir /opt/healthcheck

script
    ./run.sh
end script
