# decommission-tool-deb/src/deb/control/postinst
#!/bin/sh

adduser \
  --quiet \
  --system \
  --home /opt/decommission-tool \
  --no-create-home \
  --shell /bin/bash \
  --group \
  --gecos "MyGov Site Decommissioning Tool" \
  decommission

mkdir -p /var/log/decommission
chmod 755 /var/log/decommission
chown decommission:nogroup /var/log/decommission

systemctl daemon-reload
systemctl enable decommission
systemctl start decommission
