# decommission-tool-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop decommission || true
systemctl disable decommission || true
