# decommission-tool-deb/src/deb/control/prerm
#!/bin/sh
if [ -x /bin/systemctl ]
then
    systemctl stop decommission  || true
    systemctl disable decommission  || true
else
    service decomission stop || true
fi
